# skeleton/etc/systemd/system/fr24feed.service.d/fix-config.conf
[Service]
ExecStartPre=-/usr/local/bin/fix-config.sh
