# skeleton/etc/cron.d/weekly_reboot
# weekly reboot on Monday at between 4 AM and 5 AM
00 4 * * 1 root /bin/bash -c 'sleep $(( RANDOM \% 3600 )); /sbin/reboot;'
