# skeleton/etc/systemd/system/dhcpcd.service.d/restart.conf
[Service]
Restart=always
RestartSec=5s
StartLimitInterval=1
StartLimitBurst=100
